# profile.d/mail.sh
export MAIL=~/Maildir
export MAIL=~/Maildir
